# bin/xenon.service
[Unit]
Description=XENON User Management server
After=network.target

[Service]
ExecStart=/usr/bin/node /var/www/xenon-users-web/app.js
EnvironmentFile=/var/www/xenon-users-web/.env
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=xenonUM

[Install]
WantedBy=multi-user.target