# bin/xenon.service
[Unit]
Description=XENON User Management server
Requires=After=network.target

[Service]
ExecStart=/usr/bin/npm start
WorkingDirectory=/var/www/xenon-users-web
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=xenonUM
Environment=.env

[Install]
WantedBy=multi-user.target